# tools/packaging/linux/debian/etc/xdg/autostart/radioform-filter.desktop
[Desktop Entry]
Type=Application
Name=Radioform EQ
Exec=radioform-filter --preset /tmp/radioform-preset.json
NoDisplay=true
X-GNOME-Autostart-enabled=true
